# debian/tatu-repo-update.cron.daily
#!/bin/sh
# Ejecuta ceibal-update periodicamente para obtener actualizaciones
# de paquetes requeridos del Plan Ceibal.

set -e

[ -x /usr/sbin/ceibal-update ] || exit 0

/usr/sbin/ceibal-update
