# debian/tatu-repo-update.cron.d
# Ejecutar ceibal-update-check cada 10 minutos.
PATH=/sbin:/usr/sbin:/bin:/usr/bin

*/10 * * * * root /usr/sbin/ceibal-update-check
